# contrib/systemd/archivas@.service
[Unit]
Description=Archivas Node %i
Documentation=https://github.com/ArchivasNetwork/archivas
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=archivas
Group=archivas
WorkingDirectory=/opt/archivas

# Main service
ExecStart=/opt/archivas/bin/archivas-node \
  --rpc :8080 \
  --p2p :9090 \
  --db /opt/archivas/data/%i \
  --genesis /opt/archivas/genesis/devnet.genesis.json \
  --network-id archivas-devnet-v3 \
  --bootnodes 57.129.148.132:9090,72.251.11.191:9090 \
  --enable-gossip \
  --max-peers 20

Restart=always
RestartSec=5
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=archivas-%i

[Install]
WantedBy=multi-user.target

