# contrib/systemd/archivas-timelord.service
[Unit]
Description=Archivas Timelord (VDF Computer)
Documentation=https://github.com/ArchivasNetwork/archivas
After=archivas-node.service
Requires=archivas-node.service

[Service]
Type=simple
User=archivas
Group=archivas
WorkingDirectory=/opt/archivas

ExecStart=/opt/archivas/archivas-timelord \
  --node http://localhost:8080 \
  --step 500 \
  --metrics-addr 0.0.0.0:9101

# v1.1.1: Metrics server on port 9101 (prometheus scrape target)

Restart=always
RestartSec=5

# Logging
StandardOutput=append:/var/log/archivas/timelord.log
StandardError=append:/var/log/archivas/timelord-error.log

[Install]
WantedBy=multi-user.target

