# deploy/betanet/archivas-betanet.service
[Unit]
Description=Archivas Betanet Node
Documentation=https://docs.archivas.ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=archivas
Group=archivas

# Working directory
WorkingDirectory=/opt/archivas

# Environment
Environment="NETWORK=betanet"
Environment="DATA_DIR=/var/lib/archivas/betanet"
Environment="GOMAXPROCS=4"

# Command
ExecStart=/usr/local/bin/archivas-node \
    --network betanet \
    --rpc 0.0.0.0:8545 \
    --p2p 0.0.0.0:30303 \
    --db /var/lib/archivas/betanet \
    --genesis /etc/archivas/betanet/genesis-betanet.json \
    --max-peers 50 \
    --enable-gossip \
    --peer-whitelist seed1.betanet.archivas.ai:30303 \
    --peer-whitelist seed2.betanet.archivas.ai:30303

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitInterval=300

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=archivas-betanet

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/archivas/betanet
ReadOnlyPaths=/etc/archivas/betanet

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

